# deployment/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; FlareSolverr — handles Cloudflare bypass (headless Chrome)
[program:flaresolverr]
command=/opt/flaresolverr/flaresolverr
environment=LOG_LEVEL="info",PORT="8191",HOST="0.0.0.0",DISPLAY=":99"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Virtual display for Chrome (needed in headless Linux)
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x720x24
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Kiwi Stream API — waits 8s for FlareSolverr to start
[program:api]
command=bash -c "sleep 8 && uvicorn api_server:app --host 0.0.0.0 --port 8000"
directory=/app
environment=FLARESOLVERR_URL="http://localhost:8191/v1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
