# deployment/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; FlareSolverr — handles Cloudflare bypass
[program:flaresolverr]
command=python -m flaresolverr
environment=LOG_LEVEL="info",PORT="8191"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Kiwi Stream API — waits 5s for FlareSolverr to start
[program:api]
command=bash -c "sleep 5 && uvicorn api_server:app --host 0.0.0.0 --port 8000"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
